# etc/profile.d/alias.sh
# /etc/profile.d/alias.sh

alias df="df -h"
alias ls="ls -hSF"
alias mkdir="mkdir -p"
alias mv="mv -vi"
alias rm="rm -vi"
