# etc/profile.d/alias.sh
# /etc/profile.d/alias.sh

alias df="df -h"
alias free="free -h"
alias ls="-hSF"
alias mkdir="mkdir -vp"
alias mv="mv -vi"
alias rm="rm -vi"
