# etc/profile.d/rust.sh
# /etc/profile.d/rust.sh

export PATH="$PATH:$HOME/.cargo/bin"
